# ansible/roles/dock_ami_lc/templates/upstart.conf
#!upstart
description "{{ app_name }}"
author "Jorge Silva"

env NPM_BIN=/usr/local/bin/npm
env APP_DIR=/opt/runnable/{{ app_name }}
env LOG_FILE={{ app_log_dir }}/{{ app_name }}.log
env NODE_ENV={{ node_env }}

{% if enviroment_vars is defined %}
{% for name, value in enviroment_vars.iteritems() %}
env {{ name }}={{ value }}
{% endfor %}
{% endif %}

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

script
  touch $LOG_FILE
  chdir $APP_DIR
  echo $$ > /var/run/{{ app_name }}.pid
  exec $NPM_BIN start >> $LOG_FILE 2>&1
end script

pre-start script
  # Date format same as (new Date()).toISOString() for consistency
  echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> $LOG_FILE
end script

pre-stop script
  rm /var/run/{{ app_name }}.pid
  echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> $LOG_FILE
end script

post-start script
  echo "===== App restarted =====" >> $LOG_FILE
end script

respawn
respawn limit 5 1     # give up restart after 5 respawns in 1 seconds
